# blackout-notify/rootfs/etc/services.d/blackout-notify/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Blackout Notify
# Power grid monitoring with Telegram notifications
# ==============================================================================

# Read configuration from Home Assistant
export TELEGRAM_TOKEN=$(bashio::config 'telegram_token')
export ALLOWED_CHAT_IDS=$(bashio::config 'allowed_chat_ids')
export LOG_LEVEL=$(bashio::config 'log_level')
export POLLING_INTERVAL=$(bashio::config 'polling_interval')

# Power monitoring settings
export NOTIFICATION_CHAT_IDS=$(bashio::config 'notification_chat_ids')
export WATCHED_ENTITY_ID=$(bashio::config 'watched_entity_id')
export NEXT_ON_SENSOR_ID=$(bashio::config 'next_on_sensor_id')
export NEXT_OFF_SENSOR_ID=$(bashio::config 'next_off_sensor_id')
export PAUSE_ENTITY_ID=$(bashio::config 'pause_entity_id')
export TIMEZONE=$(bashio::config 'timezone')

# Home Assistant API URL and token
# SUPERVISOR_TOKEN is automatically available in add-on container
export HA_API_URL="http://supervisor/core/api"
export HA_TOKEN="${SUPERVISOR_TOKEN}"

# Log startup
bashio::log.info "Starting Blackout Notify..."
bashio::log.info "Log level: ${LOG_LEVEL}"
bashio::log.info "Watched entity: ${WATCHED_ENTITY_ID}"

# Check if token is configured
if bashio::var.is_empty "${TELEGRAM_TOKEN}"; then
    bashio::log.fatal "Telegram token is not configured!"
    bashio::exit.nok
fi

# Start the bot
exec /app/blackout-notify
